# plocate/packages/plocate/postinst
!#/bin/bash

systemctl enable plocate-updatedb.service
systemctl enable plocate-updatedb.timer
