# plocate/packages/plocate/postrm
!#/bin/bash

systemctl disable plocate-updatedb.service
systemctl disable plocate-updatedb.timer
