# plocate/packages/plocate/postinst
#!/bin/bash

systemctl preset plocate-updatedb.service
systemctl preset plocate-updatedb.timer

if systemctl is-enabled plocate-updatedb.timer &>/dev/null; then
    systemctl start --no-block plocate-updatedb.timer plocate-updatedb.service || :
fi
